# deploy/systemd/web3-intelligence.service
[Unit]
Description=Web3 Intelligence Telegram Bot
After=network.target

[Service]
Type=simple
User=olamide3121
WorkingDirectory=/home/olamide3121/web3-intelligence
EnvironmentFile=/home/olamide3121/web3-intelligence/.env
ExecStart=/home/olamide3121/web3-intelligence/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
