# deploy/systemd/web3-intelligence.service
[Unit]
Description=Web3 Intelligence Telegram Bot
After=network.target

[Service]
Type=simple
User=olamide3121
WorkingDirectory=/home/olamide3121/intelligence-bot 
EnvironmentFile=/home/olamide3121/intelligence-bot/.env
ExecStart=/home/olamide3121/intelligence-bot/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
